# chef/cookbooks/provisioner/templates/default/provisioner.service.erb
[Unit]
Description=Opencrowbar Provisioner Webserver
After=network.target

[Service]
ExecStart=/usr/local/bin/sws --listen=:<%= @port %> --site=<%= @docroot %>

[Install]
WantedBy=multi-user.target
